# debian/quantum-plugin-linuxbridge-agent.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
	update-rc.d quantum-plugin-linuxbridge-agent remove >/dev/null
fi

#DEBHELPER#
